# python/machine learning/crontab
*   *   *   *   * $HOME/programacion/cursos/python/ahora.sh
